# APP/sd/supervisord.conf
[supervisord]  
nodaemon=true
  
[program:sd]
command=php /mnt/apps/bin/start_swoole_server.php start -d